# deploy/quantoracle-backup.service
[Unit]
Description=QuantOracle metrics.db rotating backup (one-shot)
After=quantoracle.service

[Service]
Type=oneshot
ExecStart=/opt/quantoracle/backup-db.sh
