# deploy/quantoracle-sampler.timer
[Unit]
Description=Run QuantOracle live-data sampler every 5 minutes

[Timer]
OnCalendar=*-*-* *:00/5:00
AccuracySec=30s

[Install]
WantedBy=timers.target
